# profile.d/aliases.sh
# FILESYSTEM
# filesystem navigation, listing, and operating

alias ..='cd ..'
alias ...='cd ../..'
alias ....='cd ../../..'
alias .....='cd ../../../..'

alias ll='ls -l'
alias lsa='ls -A'
alias lla='ls -lA'
alias tra="tree -a -F -I '.git|.hg|.DS_Store|.sass-cache'"
alias tre='tree --dirsfirst -F'


# GIT
# git shorthands and utils

alias gs='git status'
alias tiga='tig --all'


# MISC
# miscellaneous cli shortcuts

alias ccat='pygmentize -O style=monokai -f console256 -g'

